# images/apache/etc/s6-overlay/s6-rc.d/apache2/run
#!/command/execlineb -P
with-contenv

env

# Set healthcheck sleep delay to 20ms
s6-notifyoncheck -s 20
/usr/sbin/apache2ctl -DFOREGROUND